# harmonista.service
[Unit]
Description=Harmonista Web Application
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/harmonista
Environment="PATH=/usr/local/go/bin:/usr/bin:/bin"
Environment="GIN_MODE=release"
EnvironmentFile=/var/www/harmonista/.env
ExecStart=/var/www/harmonista/harmonista
Restart=always
RestartSec=5

# Configurações de segurança
# NoNewPrivileges precisa estar desabilitado para usar AmbientCapabilities
# NoNewPrivileges=true
PrivateTmp=true

# Permite ao processo se conectar a portas privilegiadas (como 80 e 443)
AmbientCapabilities=CAP_NET_BIND_SERVICE

# Logs
StandardOutput=journal
StandardError=journal
SyslogIdentifier=harmonista

[Install]
WantedBy=multi-user.target
